# deploy/bbsoric.service
[Unit]
Description=BBS Oric (telnet OASCII, port 6502)
After=network.target

[Service]
Type=simple
# Utilisateur non privilégié éphémère créé par systemd (pas de compte dédié à
# provisionner ; plus sûr que root). Le BBS lit le contenu (/etc, lecture seule)
# et écrit les comptes dans son StateDirectory (/var/lib/bbsoric, RW).
DynamicUser=yes
# Répertoire d'état RW, créé et possédé par le DynamicUser, persistant entre
# redémarrages (autorisé malgré ProtectSystem=strict).
StateDirectory=bbsoric
# -metrics-addr : endpoint de supervision (/healthz, /metrics) en écoute LOCALE
# uniquement (127.0.0.1) — jamais exposé sur Internet. Sondé par bbsoric-monitor.timer.
ExecStart=/usr/local/bin/bbsoric -addr 0.0.0.0:6502 -content /etc/bbsoric/site.json -users /var/lib/bbsoric/users.json -max-conns 50 -max-conns-per-ip 3 -idle 5m -metrics-addr 127.0.0.1:6510
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bbsoric

# Durcissement (le BBS est exposé sur Internet)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
